# src/systemd/drpow.service
[Unit]
Description=DRPOW Node Service
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
EnvironmentFile=-/etc/default/drpow
ExecStart=/bin/bash -lc 'set -euo pipefail; : "${RUN_AS_USER:?missing RUN_AS_USER in /etc/default/drpow}"; : "${DRPOW_SRC_DIR:?missing DRPOW_SRC_DIR in /etc/default/drpow}"; runuser -u "${RUN_AS_USER}" -- bash -lc "cd \"${DRPOW_SRC_DIR}\" && exec ./scripts/start_testnet.sh"'
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=60
NoNewPrivileges=true
LimitNOFILE=65535

# Journald logging
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
